# conf/supervisord.d/nginx.ini
[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true      ; 在 supervisord 启动的时候也自动启动
autorestart=true    ; 程序异常退出后自动重启
startsecs=5         ; 启动 5 秒后没有异常退出，就当作已经正常启动了
stopasgroup=true    ; 通过这两项配置来确保所有子进程都能正确停止
killasgroup=true
